# 70-usb-printer.rules
DRIVER=="usblp", RUN+="/etc/init.d/cups restart"
